# contrib/systemd/srv.service
[Unit]
Description=srv Firecracker control plane
Wants=network-online.target
After=network-online.target
Requires=srv-net-helper.service srv-vm-runner.service
After=srv-net-helper.service srv-vm-runner.service

[Service]
Type=simple
User=srv
Group=srv
EnvironmentFile=/etc/srv/srv.env
ExecStart=/usr/local/bin/srv
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictRealtime=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed

[Install]
WantedBy=multi-user.target
