# images/arch-base/overlay/etc/systemd/system/srv-bootstrap.service
[Unit]
Description=srv guest tailscale bootstrap
After=local-fs.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/srv/bootstrap.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
